# trave.service
[Unit]
Description=Trave – local ML inference stack
Documentation=https://github.com/user/trave
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/oskar/projects/trave
ExecStart=/usr/bin/docker compose up -d --wait
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
